# voice_services/wyoming-satellite.service.example
[Unit]
Description=Wyoming Satellite for ReSpeaker Lite USB
After=network-online.target sound.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=lucas
WorkingDirectory=/home/lucas/ha-command-bridge/voice_services
Environment=VOICE_SATELLITE_ENV_FILE=/home/lucas/ha-command-bridge/voice_services/respeaker_lite_satellite.env
ExecStart=/home/lucas/ha-command-bridge/voice_services/run_wyoming_satellite.sh
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
